# systemd/net-watch.service
[Unit]
Description=Proxmox mgmt network watchdog (safe)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/net-reboot-if-down.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
